# apps/spotifyvm/vm-resources/local-xvfb.service
[Unit]
Description=Spotify box - Xvfb

[Service]
ExecStart=/usr/bin/Xvfb :10
User=spotify

[Install]
WantedBy=multi-user.target
